# crontab
0 0 */2 * * scrapy list | xargs -n 1 scrapy crawl
